# roles/haproxy/templates/mysqlchk.socket.j2
[Unit]
Description=Mysqlchk LibreNMS Agent Socket
 
[Socket]
ListenStream={{ mysqlchk_port }}
Accept=yes
 
[Install]
WantedBy=sockets.target
