# roles/haproxy/templates/mysqlchk@.service.j2
[Unit]
Description=Mysqlchk Agent Service
After=network.target
 
[Service]
ExecStart={{ mysqlchk_path }}
User=root
StandardOutput=socket
